# docker/root/etc/services.d/kaizoku/run
#!/usr/bin/with-contenv bash

cd /app

shopt -s globstar

umask 022

exec 2>&1
exec s6-setuidgid abc /usr/bin/node dist/server/index.js